# init.d.php-pgg.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:          php-pgg
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts php-pgg
# Description:       starts the PHP FastCGI Process Manager daemon
### END INIT INFO

prefix=@prefix@
exec_prefix=@exec_prefix@

PHP_PGG_BIN=@sbindir@/php-pgg
PHP_PGG_CONF=@sysconfdir@/php-pgg.conf
PHP_PGG_PID=@localstatedir@/run/php-pgg.pid


php_opts="--pgg-config $PHP_PGG_CONF"


wait_for_pid () {
	try=0

	while test $try -lt 35 ; do

		case "$1" in
			'created')
			if [ -f "$2" ] ; then
				try=''
				break
			fi
			;;

			'removed')
			if [ ! -f "$2" ] ; then
				try=''
				break
			fi
			;;
		esac

		echo -n .
		try=`expr $try + 1`
		sleep 1

	done

}

case "$1" in
	start)
		echo -n "Starting php-pgg "

		$PHP_PGG_BIN $php_opts

		if [ "$?" != 0 ] ; then
			echo " failed"
			exit 1
		fi

		wait_for_pid created $PHP_PGG_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	stop)
		echo -n "Gracefully shutting down php-pgg "

		if [ ! -r $PHP_PGG_PID ] ; then
			echo "warning, no pid file found - php-pgg is not running ?"
			exit 1
		fi

		kill -QUIT `cat $PHP_PGG_PID`

		wait_for_pid removed $PHP_PGG_PID

		if [ -n "$try" ] ; then
			echo " failed. Use force-quit"
			exit 1
		else
			echo " done"
		fi
	;;

	force-quit)
		echo -n "Terminating php-pgg "

		if [ ! -r $PHP_PGG_PID ] ; then
			echo "warning, no pid file found - php-pgg is not running ?"
			exit 1
		fi

		kill -TERM `cat $PHP_PGG_PID`

		wait_for_pid removed $PHP_PGG_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	restart)
		$0 stop
		$0 start
	;;

	reload)

		echo -n "Reload service php-pgg "

		if [ ! -r $PHP_PGG_PID ] ; then
			echo "warning, no pid file found - php-pgg is not running ?"
			exit 1
		fi

		kill -USR2 `cat $PHP_PGG_PID`

		echo " done"
	;;

	*)
		echo "Usage: $0 {start|stop|force-quit|restart|reload}"
		exit 1
	;;

esac
